# monitor.service
[Unit]
Description=Server Resource Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=agent
WorkingDirectory=/home/agent/workspace/monitor
ExecStart=/usr/bin/node /home/agent/workspace/monitor/server.mjs
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3099
Environment=BIND=10.0.0.3

[Install]
WantedBy=multi-user.target
